# debian/postinst
#!/bin/sh
# postinst script for amadika-shop
#
# see: dh_installdeb(1)

set -e

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <postinst> `abort-remove'
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package
SERVER_USER=www-data
SERVER_GROUP=www-data
VAR_DIR=/usr/share/amadika/shop

case "$1" in
    configure)
        pip install beautifulsoup mwlib django_dynamic_fixture django_mptt pymongo
        rm -f /etc/nginx/sites-enabled/default
        locale-gen ru_RU.utf8
        find /usr/lib/amadika/shop/ -type f -name "*.pyc" -delete

        for DIR in $VAR_DIR; do
            chown -R $SERVER_USER:$SERVER_GROUP $DIR
        done
        sudo -u $SERVER_USER /usr/lib/amadika/shop/manage.py collectstatic -v0 --noinput
        sudo -u $SERVER_USER /usr/lib/amadika/shop/manage.py synccompress -v0 --force
        sudo -u $SERVER_USER /usr/lib/amadika/shop/manage.py syncdb --migrate -v0
        update-rc.d amadika_shop defaults
        /etc/init.d/amadika_shop restart
        /etc/init.d/nginx reload
    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0
